# deploy/systemd/stemdeck-minio.service
[Unit]
Description=Stemdeck local MinIO object storage
Documentation=https://min.io/docs/minio/linux/index.html
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stemdeck
Group=stemdeck
EnvironmentFile=/etc/stemdeck/minio.env
ExecStart=/usr/local/bin/minio server $MINIO_OPTS $MINIO_VOLUMES
Restart=always
RestartSec=5
LimitNOFILE=65536

ReadWritePaths=/var/lib/stemdeck-minio
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
