# deploy/systemd/stemdeck-worker.service
[Unit]
Description=Stemdeck Celery audio processing worker
Documentation=https://github.com/your-org/stemdeck
After=network-online.target rabbitmq-server.service stemdeck-web.service
Wants=network-online.target

[Service]
Type=simple
User=stemdeck
Group=stemdeck
WorkingDirectory=/opt/stemdeck
EnvironmentFile=/etc/stemdeck/stemdeck.env
ExecStart=/opt/stemdeck/.venv/bin/python -m celery -A app.tasks.celery_app worker --loglevel=info --pool=solo
Restart=always
RestartSec=5
TimeoutStopSec=120
KillSignal=SIGTERM

# Audio separation can be memory-heavy. Raise this if your VPS has room.
MemoryMax=0

ReadWritePaths=/var/lib/stemdeck /var/cache/stemdeck
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
